# scripts/rr-morpho-maintainer.service
[Unit]
Description=Rational Reserve - Morphogenetic Maintainer
Documentation=https://github.com/morphlex/rational-reserve
After=network.target

[Service]
Type=oneshot
ExecStart=/home/llc/apps/target/release/morphlex daemon maintenance --run
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rr-morpho

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/llc/apps/data /var/log
PrivateTmp=true

# Resource limits
MemoryMax=2G
CPUQuota=100%
